# etc/hotplug.d/iface/30-vpnroute
#!/bin/sh

sleep 5
ip route add table vpn default dev singtun0
